# services/amd64/DEBIAN/postinst
systemctl daemon-reload
systemctl disable riaps-deplo.service 

