# services/amd64/DEBIAN/postinst
systemctl daemon-reload
systemctl enable riaps-rm-cgroups.service
systemctl start riaps-rm-cgroups.service
systemctl enable riaps-rm-quota.service
systemctl start riaps-rm-quota.service
systemctl enable riaps-rpyc-registry.service
systemctl start riaps-rpyc-registry.service