# services/amd64/DEBIAN/postinst
systemctl daemon-reload
systemctl disable riaps-deplo.service 
systemctl enable riaps-rpyc-registry.service
systemctl start riaps-rpyc-registry.service