# services/amd64/DEBIAN/prerm

systemctl disable riaps-disco.service 
systemctl disable riaps-deplo.service 
systemctl stop riaps-deplo.service || true 
systemctl stop riaps-deplo.service || true 

