# DEBIAN/amd64/etc/profile.d/20-riaps-env.sh
for i in `cat /etc/riaps/env.conf`;do export $i ;done
