# services/armhf/DEBIAN/postinst
systemctl daemon-reload
systemctl enable riaps-deplo.service 
systemctl start riaps-deplo.service



