# services/armhf/DEBIAN/postinst
systemctl daemon-reload
systemctl disable riaps-disco.service 
systemctl disable riaps-deplo.service 




