# services/armhf/etc/profile.d/20-riaps-env.sh
touch /etc/riaps/env.conf 
for i in `cat /etc/riaps/env.conf`;do export $i ;done
