# services/armhf/etc/systemd/system/riaps-deplo.service
[Unit]
Description=RIAPS Deployment Manager Service
Requires=riaps-rm-cgroups.service riaps-rm-quota.service

[Service]
Type=simple
EnvironmentFile=/etc/riaps/systemdenv.conf
StandardOutput=syslog
StandardError=inherit
SyslogIdentifier=RIAPS-DEPLO
SyslogLevel=info
ExecStart=/usr/local/bin/riaps_deplo
Restart=always
KillMode=process
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
