# DEBIAN/arm64/pkgfiles/postinst
set -e
pip3 install /opt/riaps-pycom/src --verbose --no-index
mkdir -p /etc/riaps
mkdir -p /usr/local/riaps/keys/

mkdir -p /usr/local/riaps/etc
rm -f /usr/local/riaps/etc/riaps.conf || true
rm -f /usr/local/riaps/etc/riaps-log.conf || true
rm -f /usr/local/riaps/keys/id_rsa.key || true
rm -f /usr/local/riaps/keys/id_rsa.pub || true
rm -f /usr/local/riaps/keys/riaps-sys.cert || true
rm -f /usr/local/riaps/keys/x509.pem || true

if [ -d "/home/riaps/.ssh/riaps_keys" ]; then
  echo "Restoring saved security keys"
  sudo cp /home/riaps/.ssh/riaps_keys/id_rsa.key /etc/riaps/id_rsa.key
  sudo cp /home/riaps/.ssh/riaps_keys/id_rsa.pub /etc/riaps/id_rsa.pub
  sudo cp /home/riaps/.ssh/riaps_keys/riaps-sys.cert /etc/riaps/riaps-sys.cert
  sudo cp /home/riaps/.ssh/riaps_keys/x509.pem /etc/riaps/x509.pem
fi

chown root:riaps /etc/riaps/id_rsa.key
chown root:riaps /etc/riaps/id_rsa.pub
chown root:riaps /etc/riaps/riaps-sys.cert
chown root:riaps /etc/riaps/x509.pem
chmod 440 /etc/riaps/id_rsa.key
chmod 440 /etc/riaps/id_rsa.pub
chmod 444 /etc/riaps/riaps-sys.cert
chmod 440 /etc/riaps/x509.pem

ln -s /etc/riaps/riaps.conf /usr/local/riaps/etc/riaps.conf
ln -s /etc/riaps/riaps-log.conf /usr/local/riaps/etc/riaps-log.conf
ln -s /etc/riaps/id_rsa.key /usr/local/riaps/keys/id_rsa.key
ln -s /etc/riaps/id_rsa.pub /usr/local/riaps/keys/id_rsa.pub
ln -s /etc/riaps/riaps-sys.cert /usr/local/riaps/keys/riaps-sys.cert
ln -s /etc/riaps/x509.pem /usr/local/riaps/keys/x509.pem

#Create soft link from rdiscoveryd to run as default
mv /usr/local/bin/riaps_disco /usr/local/bin/riaps_disco_redis
ln -s /usr/local/bin/rdiscoveryd /usr/local/bin/riaps_disco

rm -rf /opt/riaps-pycom/ || true

systemctl daemon-reload
systemctl enable riaps-rm-cgroups.service
systemctl start riaps-rm-cgroups.service || true
systemctl enable riaps-rm-quota.service
systemctl start riaps-rm-quota.service || true
systemctl enable riaps-deplo.service
systemctl start riaps-deplo.service || true
systemctl enable sethostname.service
systemctl start sethostname.service || true
