# DEBIAN/arm64/pkgfiles/prerm

systemctl disable riaps-deplo.service
systemctl disable sethostname.service
systemctl disable riaps-rm-cgroups.service
systemctl disable riaps-rm-quota.service
systemctl stop riaps-deplo.service || true
systemctl stop sethostname.service || true
systemctl stop riaps-rm-cgroups.service || true
systemctl stop riaps-rm-quota.service || true
